# packaging/launchagents/ai.hermes.gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.hermes.gateway</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/mong-e/workspace/hermes-agent/venv/bin/python</string>
        <string>-m</string>
        <string>hermes_cli.main</string>
        <string>gateway</string>
        <string>run</string>
        <string>--replace</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mong-e/workspace/hermes-agent</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/mong-e/workspace/hermes-agent/venv/bin:/Users/mong-e/workspace/hermes-agent/node_modules/.bin:/opt/homebrew/bin:/opt/homebrew/sbin:/Users/mong-e/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>VIRTUAL_ENV</key>
        <string>/Users/mong-e/workspace/hermes-agent/venv</string>
        <key>HERMES_HOME</key>
        <string>/Users/mong-e/.hermes</string>
        <key>HERMES_FC_RAG_ENABLED</key>
        <string>1</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/mong-e/.hermes/logs/gateway.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/mong-e/.hermes/logs/gateway.error.log</string>
</dict>
</plist>
