# packaging/launchagents/ai.hermes.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.hermes.monitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/mong-e/.hermes/venv/bin/python3</string>
        <string>/Users/mong-e/workspace/hermes-agent/dashboard/hermes-monitor-mac.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/mong-e/.hermes/logs/monitor-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mong-e/.hermes/logs/monitor-stderr.log</string>
</dict>
</plist>
